# supervisord.conf
[unix_http_server]
file = /tmp/supervisord.sock

[supervisord]
logfile = /var/log/supervisord.log
pidfile = /var/run/supervisord.pid
loglevel = info
directory = /usr/src
childlogdir = /var/log/supervisord

[supervisorctl]
serverurl = unix:///tmp/supervisord.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[program:celery]
command = /usr/src/scripts/celery.grid.sh
environment = CROSSWATCH_CELERY_BIN="celery"
autorestart = false
